# crontab

@reboot nohup /bin/sh ~/scripts/on_reboot.sh
