# crontab

@reboot /bin/sh ~/scripts/on_reboot.sh
